# scripts/systemd/cert-renewal.service
[Unit]
Description=SSL Certificate Renewal (Tailscale + Let's Encrypt)
Documentation=https://tailscale.com/kb/1153/enabling-https/
Documentation=https://certbot.eff.org/
After=network-online.target tailscaled.service docker.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/aragorn/scripts/renew-certs.sh
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=false
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/mnt/secondary/containers/nginx/ssl /var/log /etc/letsencrypt
PrivateTmp=true

[Install]
WantedBy=multi-user.target

