# scripts/systemd/tailscale-cert-renewal.service
[Unit]
Description=Tailscale Certificate Renewal for Nginx
Documentation=https://tailscale.com/kb/1153/enabling-https/
After=network-online.target tailscaled.service
Wants=network-online.target
Requires=tailscaled.service

[Service]
Type=oneshot
ExecStart=/opt/aragorn/scripts/renew-tailscale-certs.sh
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=false
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/mnt/secondary/containers/nginx/ssl /var/log
PrivateTmp=true

# Environment
Environment="TAILSCALE_DOMAIN=casa.home"
Environment="NGINX_SSL_DIR=/mnt/secondary/containers/nginx/ssl"
Environment="COMPOSE_DIR=/opt/aragorn"
Environment="LOG_FILE=/var/log/tailscale-cert-renewal.log"
Environment="RENEWAL_DAYS=30"

[Install]
WantedBy=multi-user.target

