# scripts/systemd/tailscale-cert-renewal.timer
[Unit]
Description=Monthly Tailscale Certificate Renewal Timer
Documentation=https://tailscale.com/kb/1153/enabling-https/

[Timer]
# Esegue il 1° giorno di ogni mese alle 03:00
OnCalendar=*-*-01 03:00:00

# Esegue immediatamente se l'ultimo run è stato perso
Persistent=true

# Randomizza l'esecuzione di ±1 ora per evitare picchi
RandomizedDelaySec=3600

# Accuracy del timer
AccuracySec=1min

[Install]
WantedBy=timers.target

