# scripts/crontab
# Crontab per il rinnovo automatico dei certificati Tailscale
# 
# Installazione:
#   sudo crontab -e
#   # Incollare la riga seguente
#
# Oppure:
#   sudo crontab /path/to/aragorn/scripts/crontab
#

# Rinnovo certificati Tailscale - Esegue il 1° giorno di ogni mese alle 03:00
0 3 1 * * /opt/aragorn/scripts/renew-tailscale-certs.sh >> /var/log/tailscale-cert-renewal.log 2>&1

